# systemd/readsb.service.d/10-network-mode.conf
[Unit]
ConditionPathExists=!/opt/adsbledmatrix/data/.network_receiver_enabled
