# systemd/adsbledmatrix.service
[Unit]
Description=ADS-B LED Display
After=network.target readsb.service
Wants=readsb.service

[Service]
Type=simple
# rpi-rgb-led-matrix needs root to initialise GPIO timing hardware.
# It drops privileges after init (see hardware/led_matrix.py).
User=root
Group=root
WorkingDirectory=/opt/adsbledmatrix
Environment=PYTHONPATH=/opt/adsbledmatrix/backend:/opt/adsbledmatrix
Environment=ADSB_DEBUG=false
Environment=ADSB_DATA_DIR=/opt/adsbledmatrix/data
ExecStart=/opt/adsbledmatrix/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
