# systemd/adsbledmatrix-update.service
[Unit]
Description=ADS-B LED Matrix update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/adsbledmatrix/scripts/check_and_update.sh
